# debian/iprutils.iprinit.upstart
# iprinit - ipr initialization daemon

description     "ipr initialization daemon"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn
respawn limit 10 5

pre-start script
    if [ ! -d /sys/class/scsi_generic ] ; then
        modprobe sg
    fi
end script

exec /sbin/iprinit --daemon
